# deploy/krv-notify.timer
[Unit]
Description=Daily KRV discovery bridge (notify Мангуста о новых бесплатных ключах)

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
